# templates/mediacenter/etc/cron.daily/backup
#!/bin/bash
rsync -az --delete /mnt/ROBA01/ /mnt/ROBA02
